# Backup/readymedia-realtime/readymedia-realtime.service
[Unit]
Description=ReadyMedia Realtime - Universal Accessibility Transcription Service
Documentation=https://github.com/din-bruker/readymedia-realtime
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/readymedia-realtime/server

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/var/www/readymedia-realtime/server/.env

# Start command
ExecStart=/usr/bin/node server.js

# Restart policy
Restart=on-failure
RestartSec=10
KillMode=process

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=readymedia-realtime

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/readymedia-realtime

[Install]
WantedBy=multi-user.target
